# service1.service
[Unit]
Description=benim servisim


[service]
Restart=always
ExecStart=/home/alper/bsm/script.sh

[Install]
WantedBy=multi-user.target
